# single-node-kafka.conf
description "Single-node Kafka"
author "alex@woolfordcompany.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/kafka/kafka-0.8.1.1-src/start-single-node-kafka.sh
